# misc/clipboard_sync/clipboard-sync-host.service
[Unit]
Description=Clipboard Sync Host Service
After=graphical-session.target

[Service]
ExecStart=%h/.dotfiles/misc/clipboard_sync/clipboard_sync.sh host
Restart=always
Environment=WAYLAND_DISPLAY=%E/wayland-0
Environment=XDG_RUNTIME_DIR=%E/XDG_RUNTIME_DIR

[Install]
WantedBy=default.target

